# deploy/rambo-bot.service
[Unit]
Description=Rambo Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rambo
Group=rambo
WorkingDirectory=/opt/rambo-bot
ExecStartPre=+/bin/bash /opt/rambo-bot/deploy/load-metadata.sh
ExecStart=/bin/bash -c 'set -a; . /opt/rambo-bot/.env.runtime; set +a; exec /usr/bin/node index.js'
Restart=on-failure
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=rambo-bot

[Install]
WantedBy=multi-user.target
